# chef-repo/cookbooks/grailsapp/templates/default/grailsapp.service.erb
[Unit]
Description=grailsapp Java Daemon

[Service]
ExecStart=/usr/bin/java -jar <%= node['grailsapp']['jar'] %>
User=<%= node['grailsapp']['user'] %>

[Install]
WantedBy=multi-user.target
